# .anacron/cron.daily/test.sh
#!/bin/bash
date >> ${HOME}/Desktop/dateytest.txt
